# Rasperry Pi Discovery Tool/systemd/initial-checkin.service
[Unit]
Description=Yeyland Wutani Network Discovery Pi - Initial Check-In
Documentation=https://github.com/the-last-one-left/YeylandWutani
# Run after network is fully up AND the system clock is NTP-synced.
# time-sync.target prevents the Azure AD token lifetime-validation error
# that occurs when the Pi boots with a stale/wrong clock.
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target
# Must run before the full discovery service
Before=network-discovery.service

[Service]
Type=oneshot
User=network-discovery
Group=network-discovery
WorkingDirectory=/opt/network-discovery
ExecStartPre=+/bin/bash -c "mkdir -p /opt/network-discovery/data /opt/network-discovery/logs && chown network-discovery:network-discovery /opt/network-discovery/data /opt/network-discovery/logs && chmod 750 /opt/network-discovery/data /opt/network-discovery/logs && chown root:network-discovery /opt/network-discovery/config && chmod 750 /opt/network-discovery/config"
ExecStart=/opt/network-discovery/venv/bin/python3 /opt/network-discovery/bin/initial-checkin.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nd-checkin
TimeoutStartSec=300

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/network-discovery/logs /opt/network-discovery/data
ProtectHome=yes
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
