# Rasperry Pi Discovery Tool/systemd/network-discovery.service
[Unit]
Description=Yeyland Wutani Network Discovery Pi - Full Network Discovery
Documentation=https://github.com/the-last-one-left/YeylandWutani
After=network-online.target initial-checkin.service
Wants=network-online.target
# Wait for initial check-in to succeed first
Requires=initial-checkin.service

[Service]
Type=simple
User=network-discovery
Group=network-discovery
WorkingDirectory=/opt/network-discovery
ExecStartPre=+/bin/bash -c "mkdir -p /opt/network-discovery/data /opt/network-discovery/logs && chown network-discovery:network-discovery /opt/network-discovery/data /opt/network-discovery/logs && chmod 750 /opt/network-discovery/data /opt/network-discovery/logs && chown root:network-discovery /opt/network-discovery/config && chmod 750 /opt/network-discovery/config"
ExecStart=/opt/network-discovery/venv/bin/python3 /opt/network-discovery/bin/discovery-main.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nd-discovery
TimeoutStartSec=1800
TimeoutStopSec=60
Restart=no

# Grant raw socket access for nmap/arp-scan
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/network-discovery/logs /opt/network-discovery/data
ProtectHome=yes

[Install]
WantedBy=multi-user.target
