# Rasperry Pi Risk Scanner Tool/systemd/risk-scanner-daily.timer
[Unit]
Description=Yeyland Wutani Risk Scanner - Daily Scan Timer
Requires=risk-scanner-daily.service

[Timer]
OnCalendar=*-*-* __SCAN_TIME__:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
