# Rasperry Pi Risk Scanner Tool/systemd/risk-scanner-report.timer
[Unit]
Description=Yeyland Wutani Risk Scanner - Weekly Report Timer
Requires=risk-scanner-report.service

[Timer]
OnCalendar=__REPORT_DAY__ __REPORT_TIME__:00

[Install]
WantedBy=timers.target
